# debian/metaproxy.prerm
#!/bin/sh
if [ -x "/etc/init.d/metaproxy" ]; then
	invoke-rc.d metaproxy stop || exit $?
fi
